# chroot/etc/profile.d/man.sh
#/bin/sh
test -f /usr/share/man/mandoc.db || sudo makewhatis
