# chroot/etc/profile.d/set_vi.sh
#!/bin/sh
set -o vi
export EDITOR=vim
export FCEDIT=$EDITOR
